# roles/application/files/application.conf
[program:application]
command=python /opt/application/bootstrap.py 8000
autostart=true
autorestart=true
stderr_logfile=/var/log/application.err.log
stdout_logfile=/var/log/application.out.log